# deploy/demos2sales.service
[Unit]
Description=Demos2Sales React + FastAPI app
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/demos2sales/current/backend
EnvironmentFile=/opt/demos2sales/shared/.env
Environment=DEMO_CALC_DB_HOST=127.0.0.1
Environment=DEMO_CALC_DB_PORT=5464
Environment=DEMO_CALC_DB_NAME=demo_calc
Environment=DEMO_CALC_DB_USER=admin
Environment=DEMO_CALC_AUTO_SEED_EMPTY_DB=0
Environment=DEMO_CALC_JSON_BACKUP_PATH=/opt/demos2sales/shared/data/state_backup.json
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/demos2sales/current/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 7861
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
